# files/ubuntu-root/etc/s6-overlay/s6-rc.d/svc-auth/run
#!/usr/bin/with-contenv bash
exec /usr/bin/python3 /usr/local/bin/auth-server.py
