# files/ubuntu-root/etc/s6-overlay/s6-rc.d/svc-selkies/run
#!/usr/bin/with-contenv bash

TARGET_USER="${CUSTOM_USER:-${USER_NAME:-root}}"
PORT="${SELKIES_PORT:-8082}"

# ENCODER is passed from start-container.sh: software|nvidia|nvidia-wsl|intel|amd
ENCODER="${ENCODER:-}"
GPU_VENDOR="${ENCODER:-${GPU_VENDOR:-software}}"

# Check if running in WSL environment
WSL_ENVIRONMENT="${WSL_ENVIRONMENT:-false}"

# Export LIBVA_DRIVER_NAME for Intel/AMD VA-API (set by start-container.sh)
if [ -n "${LIBVA_DRIVER_NAME}" ]; then
  export LIBVA_DRIVER_NAME
  echo "[svc-selkies] Using VA-API driver: ${LIBVA_DRIVER_NAME}"
fi

# WSL2 specific setup for NVIDIA - add WSL libs to paths
if [ "${WSL_ENVIRONMENT}" = "true" ]; then
  echo "[svc-selkies] Running in WSL2 environment"
  if [ -d "/usr/lib/wsl/lib" ]; then
    # Remove /usr/lib/wsl/lib if already present, then add to front
    LD_LIBRARY_PATH="${LD_LIBRARY_PATH//\/usr\/lib\/wsl\/lib:/}"
    LD_LIBRARY_PATH="${LD_LIBRARY_PATH//\/usr\/lib\/wsl\/lib/}"
    export LD_LIBRARY_PATH="/usr/lib/wsl/lib:${LD_LIBRARY_PATH}"
    export PATH="/usr/lib/wsl/lib:${PATH}"
    echo "[svc-selkies] Added WSL library path at front: /usr/lib/wsl/lib"
  fi
fi

# Function to find nvidia-smi (handles WSL2 where it's in /usr/lib/wsl/lib)
find_nvidia_smi() {
  if command -v nvidia-smi &> /dev/null; then
    echo "nvidia-smi"
  elif [ -x "/usr/lib/wsl/lib/nvidia-smi" ]; then
    echo "/usr/lib/wsl/lib/nvidia-smi"
  else
    return 1
  fi
}

# Check if NVIDIA GPU is available (works on both native Linux and WSL2)
nvidia_available() {
  local nvidia_smi_cmd
  nvidia_smi_cmd=$(find_nvidia_smi) || return 1
  "${nvidia_smi_cmd}" >/dev/null 2>&1
}

# Configure pixelflux hardware encoding based on GPU_VENDOR
# pixelflux supports: NVENC (NVIDIA), VA-API (Intel/AMD)
# Setting vaapi_render_node_index >= 0 enables VA-API encoding on /dev/dri/renderD(128 + index)
# -1 means try NVENC if available, otherwise software

# Determine GPU render node index for pixelflux
if [ "${GPU_VENDOR}" = "nvidia" ] || [ "${GPU_VENDOR}" = "nvidia-wsl" ]; then
  # NVIDIA: pixelflux will auto-detect NVENC, set to -1
  export PIXELFLUX_VAAPI_RENDER_NODE_INDEX="-1"
  echo "[svc-selkies] GPU_VENDOR=${GPU_VENDOR}: pixelflux will use NVENC hardware encoding"
  
  if nvidia_available; then
    NVIDIA_SMI_CMD=$(find_nvidia_smi)
    echo "[svc-selkies] NVIDIA GPU detected via: ${NVIDIA_SMI_CMD}"
  fi
elif [ "${GPU_VENDOR}" = "intel" ] || [ "${GPU_VENDOR}" = "amd" ]; then
  # Intel/AMD: Use VA-API with render node
  # Need to find the correct render node for the specified GPU vendor
  VAAPI_INDEX=""
  
  # Find the correct render node by checking the driver
  for renderD in /dev/dri/renderD*; do
    if [ -e "$renderD" ]; then
      RENDER_NUM=$(echo "$renderD" | grep -oP 'renderD\K[0-9]+')
      DRIVER_LINK="/sys/class/drm/renderD${RENDER_NUM}/device/driver"
      if [ -L "$DRIVER_LINK" ]; then
        DRIVER_NAME=$(basename "$(readlink -f "$DRIVER_LINK")")
        if [ "${GPU_VENDOR}" = "intel" ] && [ "$DRIVER_NAME" = "i915" ]; then
          VAAPI_INDEX=$((RENDER_NUM - 128))
          echo "[svc-selkies] Found Intel GPU at renderD${RENDER_NUM} (driver: ${DRIVER_NAME})"
          break
        elif [ "${GPU_VENDOR}" = "amd" ] && [ "$DRIVER_NAME" = "amdgpu" ]; then
          VAAPI_INDEX=$((RENDER_NUM - 128))
          echo "[svc-selkies] Found AMD GPU at renderD${RENDER_NUM} (driver: ${DRIVER_NAME})"
          break
        fi
      fi
    fi
  done
  
  if [ -n "${VAAPI_INDEX}" ]; then
    export PIXELFLUX_VAAPI_RENDER_NODE_INDEX="${VAAPI_INDEX}"
    echo "[svc-selkies] GPU_VENDOR=${GPU_VENDOR}: pixelflux will use VA-API hardware encoding (index ${VAAPI_INDEX})"
  elif [ -n "${DRI_NODE}" ]; then
    # Fallback to DRI_NODE if set
    RENDER_NUM=$(echo "${DRI_NODE}" | grep -oP 'renderD\K[0-9]+')
    if [ -n "${RENDER_NUM}" ]; then
      VAAPI_INDEX=$((RENDER_NUM - 128))
      export PIXELFLUX_VAAPI_RENDER_NODE_INDEX="${VAAPI_INDEX}"
      echo "[svc-selkies] GPU_VENDOR=${GPU_VENDOR}: pixelflux will use VA-API on ${DRI_NODE} (index ${VAAPI_INDEX})"
    else
      export PIXELFLUX_VAAPI_RENDER_NODE_INDEX="-1"
      echo "[svc-selkies] GPU_VENDOR=${GPU_VENDOR}: Could not determine render node, falling back to software encoding"
    fi
  else
    export PIXELFLUX_VAAPI_RENDER_NODE_INDEX="-1"
    echo "[svc-selkies] GPU_VENDOR=${GPU_VENDOR}: No matching render node found, falling back to software encoding"
  fi
else
  # No GPU or unknown: software encoding
  export PIXELFLUX_VAAPI_RENDER_NODE_INDEX="-1"
  echo "[svc-selkies] GPU_VENDOR=${GPU_VENDOR}: pixelflux will use software encoding"
fi

# Always use selkies with pixelflux (WebSocket mode)
SELKIES_CMD="selkies"
CMD_OPTS=(--addr="localhost" --mode="websockets" --port="${PORT}")

# Configure encoder based on GPU vendor
ENCODER_TYPE="x264enc"  # Default to software encoder

if [ "${GPU_VENDOR}" = "nvidia" ] || [ "${GPU_VENDOR}" = "nvidia-wsl" ]; then
  ENCODER_TYPE="nvh264enc"
  CMD_OPTS+=(--encoder="${ENCODER_TYPE}")
  echo "[svc-selkies] Using NVIDIA NVENC encoder: ${ENCODER_TYPE}"
elif [ "${GPU_VENDOR}" = "intel" ] || [ "${GPU_VENDOR}" = "amd" ]; then
  # VA-API hardware encoder
  ENCODER_TYPE="vah264enc"
  CMD_OPTS+=(--encoder="${ENCODER_TYPE}")
  
  if [ -n "${VAAPI_INDEX}" ]; then
    DRI_NODE_PATH="/dev/dri/renderD$((128 + VAAPI_INDEX))"
    if [ -e "${DRI_NODE_PATH}" ]; then
      CMD_OPTS+=(--dri-node="${DRI_NODE_PATH}")
      echo "[svc-selkies] Added --dri-node=${DRI_NODE_PATH} for VA-API encoding"
    fi
  fi
  echo "[svc-selkies] Using VA-API encoder: ${ENCODER_TYPE}"
else
  # Software encoder (x264)
  CMD_OPTS+=(--encoder="${ENCODER_TYPE}")
  echo "[svc-selkies] Using software encoder: ${ENCODER_TYPE}"
fi

# Enable H.264 streaming mode for lower latency (crucial for responsiveness!)
# This enables low-latency encoding with reduced buffering
CMD_OPTS+=(--h264-streaming-mode=true)
echo "[svc-selkies] H.264 streaming mode enabled for low latency"

# Optimize CRF for better quality/latency balance (lower = better quality, more bandwidth)
# Default is 25, we use 18-20 for better responsiveness
SELKIES_H264_CRF="${SELKIES_H264_CRF:-18}"
CMD_OPTS+=(--h264-crf="${SELKIES_H264_CRF}")
echo "[svc-selkies] H.264 CRF set to: ${SELKIES_H264_CRF}"

echo "[svc-selkies] Using selkies with pixelflux (GPU_VENDOR=${GPU_VENDOR}, Encoder=${ENCODER_TYPE})"

# Default sink setup (wait until pulseaudio is ready, then create virtual sinks)
if [ ! -f '/dev/shm/audio.lock' ]; then
  # ensure pulseaudio is up
  for i in $(seq 1 30); do
    if s6-setuidgid "$TARGET_USER" with-contenv pactl info >/dev/null 2>&1; then
      READY=1
      break
    fi
    sleep 0.5
  done
  if [ "${READY:-0}" -eq 1 ]; then
    s6-setuidgid "$TARGET_USER" with-contenv pactl \
      load-module module-null-sink \
      sink_name="output" \
      sink_properties=device.description="output"
    s6-setuidgid "$TARGET_USER" with-contenv pactl \
      load-module module-null-sink \
      sink_name="input" \
      sink_properties=device.description="input"
    touch /dev/shm/audio.lock
  else
    echo "[svc-selkies] pulseaudio not ready; skipped null-sink setup (audio may be missing)." >&2
  fi
fi

# Setup dev mode if defined
if [ ! -z ${DEV_MODE+x} ]; then
  # Dev deps
  apt-get update
  apt-get install -y \
    nodejs
  npm install -g nodemon
  rm -Rf $HOME/.npm
  # Frontend setup
  if [[ "${DEV_MODE}" == "core" ]]; then
    # Core just runs from directory
    cd $HOME/src/addons/gst-web-core
    s6-setuidgid "$TARGET_USER" npm install
    s6-setuidgid "$TARGET_USER" npm run serve &
  else
    # Build core
    cd $HOME/src/addons/gst-web-core
    s6-setuidgid "$TARGET_USER" npm install
    s6-setuidgid "$TARGET_USER" npm run build
    s6-setuidgid "$TARGET_USER" cp dist/selkies-core.js ../${DEV_MODE}/src/
    s6-setuidgid "$TARGET_USER" nodemon --watch selkies-core.js --exec "npm run build && cp dist/selkies-core.js ../${DEV_MODE}/src/" & 
    # Copy touch gamepad
    s6-setuidgid "$TARGET_USER" cp ../universal-touch-gamepad/universalTouchGamepad.js ../${DEV_MODE}/src/
    s6-setuidgid "$TARGET_USER" nodemon --watch ../universal-touch-gamepad/universalTouchGamepad.js --exec "cp ../universal-touch-gamepad/universalTouchGamepad.js ../${DEV_MODE}/src/" &  
    # Copy themes
    s6-setuidgid "$TARGET_USER" cp -a nginx ../${DEV_MODE}/
    # Run passed frontend
    cd $HOME/src/addons/${DEV_MODE}
    s6-setuidgid "$TARGET_USER" npm install
    s6-setuidgid "$TARGET_USER" npm run serve &
  fi
  # Run backend
  cd $HOME/src/src
  s6-setuidgid "$TARGET_USER" \
    nodemon -V --ext py --exec \
      "python3" -m selkies \
      --addr="localhost" \
      --mode="websockets" \
      --debug="true"
fi

# Start Selkies
echo "[svc-selkies] Starting ${SELKIES_CMD} on port ${PORT}"
exec s6-setuidgid "$TARGET_USER" "${SELKIES_CMD}" "${CMD_OPTS[@]}"
