# files/ubuntu-root/etc/s6-overlay/s6-rc.d/init-selkies-end/up
# This file doesn't do anything, it's just the end of the selkies init process
