# files/ubuntu-root/etc/s6-overlay/s6-rc.d/init-selkies/up
# This file doesn't do anything, it's just the beginning of the selkies init process
